# deskflow-auto-allow.service
[Unit]
Description=Auto-accept deskflow input capture dialog
After=graphical-session.target
BindsTo=graphical-session.target
Wants=ydotool.service

[Service]
Type=simple
ExecStartPre=/bin/sleep 5
ExecStart=%h/.local/bin/deskflow-auto-allow
Restart=no

[Install]
WantedBy=default.target
