# deskflow-auto-allow.service
[Unit]
Description=Auto-accept deskflow input capture dialog
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/deskflow-auto-allow
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
