# deskflow-auto-allow.service
[Unit]
Description=Auto-accept deskflow input capture dialog
After=graphical-session.target
BindsTo=graphical-session.target

[Service]
Type=simple
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority
ExecStartPre=/bin/sleep 10
ExecStart=%h/.local/bin/deskflow-auto-allow
Restart=no

[Install]
WantedBy=default.target
